# crontab
5,35 * * * * /home/yagnesh/git/iitmtv/makegifs.sh
