# crontab
5,35 * * * * /home/yagnesh/git/monsoon_portal/makegifs.sh
